# systemd/automation.target
[Unit]
Description=Automation Job Suite Orchestration
StopWhenUnneeded=yes

[Install]
WantedBy=multi-user.target
