# systemd/auto-fetch-api-data.service
[Unit]
Description=YouTube Search Metadata Job
After=network.target
PartOf=automation.target

[Service]
Type=oneshot
User=automation-user
Group=automation-user
WorkingDirectory=/opt/automation
# Punta al nuovo script che abbiamo creato
ExecStart=/opt/automation/bin/run-youtube-metadata.sh
# Assicurati che l'ambiente abbia uv nel PATH
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/automation-user/.local/bin

[Install]
WantedBy=automation.target
